# service-configs/meetingnotes-webhook.service
[Unit]
Description=Meeting Notes Webhook Daemon
Documentation=https://github.com/ewilderj/meeting-notes-processor
After=network.target

[Service]
Type=simple

# CHANGEME: Set your username
User=CHANGEME

# CHANGEME: Path to your meeting-notes-processor repo
WorkingDirectory=/home/CHANGEME/git/meeting-notes-processor

# CHANGEME: Path to uv (find with: which uv)
# CHANGEME: Set your GitHub token
ExecStart=/home/CHANGEME/.local/bin/uv run meetingnotesd.py
Environment="GH_TOKEN=ghp_CHANGEME"

# Prevent git from prompting for credentials (fail fast instead of hanging)
Environment="GIT_ASKPASS=/bin/true"
Environment="GIT_TERMINAL_PROMPT=0"

# Optional: load secrets from a separate file instead of embedding in this unit
# EnvironmentFile=/etc/meetingnotes-webhook.env

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging (goes to journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meetingnotes-webhook

# Note: Security hardening options (NoNewPrivileges, ProtectSystem, ProtectHome, etc.)
# are intentionally omitted. They conflict with snap-based uv and npx which need
# to set capabilities and access various system paths. If you install uv via
# the standalone installer (curl -LsSf https://astral.sh/uv/install.sh | sh)
# instead of snap, you may be able to re-enable some hardening options.

[Install]
WantedBy=multi-user.target
